# templates/etc/systemd/system/docker.service.d/options.conf.j2
# {{ ansible_managed }}

[Unit]
Description=dockerdaemon for docker service
StartLimitIntervalSec=0 
#StartLimitBurst=

[Service]
ExecStart=
ExecStart=/usr/bin/dockerd
